# Tools/crontab
# Debug logs collector
* * * * * /var/spool/watcher &> /var/log/crontab
# Sync data
5 * * * * /bin/syncdata &> /dev/null
