# Tools/crontab
# Debug logs collector
* * * * * /var/spool/watcher &> /var/log/crontab
